# cron/dockermi.timer
[Unit]
Description=Timer for %p service
Requires=docker.service
[Timer]
OnCalendar=daily
[X-Fleet]
MachineOf=%p.service
