# cron/reset-failed.timer
[Unit]
Description=Timer for %p service
[Timer]
OnCalendar=weekly
[X-Fleet]
MachineOf=%p.service
